# scripts/launchd/com.switchboard.archive-superpowers.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Daily launchd job that archives superpowers specs/plans older than 2 days.
  Installed by scripts/launchd/install-archive-superpowers.sh, which substitutes
  the {{...}} placeholders below. Do not edit the installed copy by hand —
  re-run the installer instead.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.switchboard.archive-superpowers</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>cd "{{REPO_ROOT}}" &amp;&amp; ./scripts/archive-superpowers-docs.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{REPO_ROOT}}</string>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>{{HOME}}/Library/Logs/switchboard-archive.log</string>

    <key>StandardErrorPath</key>
    <string>{{HOME}}/Library/Logs/switchboard-archive.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
  </dict>
</plist>
